# kindle-connector.service
[Unit]
Description=kindle-connector — no-JS web terminal for Kindle -> claude/codex
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/codebase/kindle-connector
ExecStart=/usr/bin/python3 %h/codebase/kindle-connector/server.py
Environment=KINDLE_PORT=8090
Environment=KINDLE_COLS=58
Environment=KINDLE_ROWS=32
Environment=KINDLE_CMD=claude
# Environment=KINDLE_TOKEN=changeme
# Environment=KINDLE_WORKDIR=%h/codebase
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
